# themes/luci-theme-teleofis/root/etc/uci-defaults/35_luci-theme-teleofis
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Teleofis=/luci-static/teleofis
	set luci.main.mediaurlbase=/luci-static/teleofis
	commit luci
EOF
exit 0
